# target/linux/ipq807x/base-files/etc/init.d/nat_loopback
#!/bin/sh /etc/rc.common

START=99

fix_nat_loopback_ecm(){
    sysctl -w net.bridge.bridge-nf-call-arptables=0
    sysctl -w net.bridge.bridge-nf-call-iptables=0
    sysctl -w net.bridge.bridge-nf-call-ip6tables=0
}

boot() {
    case $(board_name) in
        dynalink,dl-wrx36|\
        edgecore,eap102|\
        edimax,cax1800|\
        qnap,301w|\
        redmi,ax6|\
        xiaomi,ax3600|\
        xiaomi,ax9000)
                fix_nat_loopback_ecm
        ;;
    esac
}
